# upstart/reddit-consumers-restart.conf
description "restart queue consumers"

task

start on reddit-restart or reddit-kill

script
    . /etc/default/reddit
    manage-consumers
end script
