# upstart/reddit-boot.conf
description "start reddit on boot"

task
start on runlevel [2345]

exec initctl emit reddit-start
