# upstart/reddit-consumers-start.conf
description "start up queue consumers"

task

start on runlevel [2345] or reddit-start

script
    . /etc/default/reddit
    manage-consumers
end script
